# ld.so.preload
/usr/lib/libfake_rlimit.so
/usr/lib/libhardened_malloc-light.so
